# setup/stock/gpionext.service
[Unit]
Description=GPIOnext GPIO to HID daemon
After=network.target

[Service]
Type=simple
User=root
# WorkingDirectory updated from /home/pi/gpionext to /opt/gpionext
WorkingDirectory=/opt/gpionext

ExecStart=/opt/gpionext/bin/gpionext --config /opt/gpionext/config/gpionext.json --debounce 1 --combo_delay 50

# SIGHUP triggers hot-reload (re-reads gpionext.json without restarting).
# Used by 'gpionext reload' and 'gpionext set *' commands.
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGTERM

# Restart automatically on crash (e.g. GPIO module not loaded)
Restart=on-failure
RestartSec=2

# Capture stdout/stderr in journald
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
